# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=service php5-fpm start
autorestart=true